# app/crontab
# Run the script every 24 hours
0 0 * * * python /app/app.py $MOUSOUTRADE_CONFIG_FILE >> /var/log/cron.log 2>&1
