# app/crontab
# All times are in America/Los_Angeles timezone (PST/PDT)
# Run the script at 2 AM PST/PDT on business days (Monday-Friday)
0 3 * * 1-5 root /bin/bash /app/startup.sh cron
