# app/crontab
# Run the script at 2 AM PST on business days (Monday-Friday)
# Redirect both stdout (1) and stderr (2) to the log file
0 2 * * 1-5 /bin/bash -c 'source /app/venv/bin/activate && python /app/run.py 1>> /var/log/cron.log 2>> /var/log/cron.error.log'
