# app/crontab
# All times are in America/Los_Angeles timezone (PST/PDT)
# Run the script at 3:06 AM PST/PDT on business days (Monday-Friday)
0 4 * * 1-5 /app/startup.sh cron
